# configs/claude/com.teazyou.patch-watcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<!-- EXPERIMENTAL: fires scripts/patch-watcher.sh every 30 min to re-apply the
	     Obsidian transparency patch and the Claude Code VS Code panel patch after
	     their respective app/extension updates wipe them. NOT wired into the
	     fresh-Mac install flow on purpose (see the script header + doc for
	     rationale and disable instructions). Symlinked into ~/Library/LaunchAgents/.

	     load:    launchctl bootstrap gui/$UID ~/Library/LaunchAgents/com.teazyou.patch-watcher.plist
	     unload:  launchctl bootout   gui/$UID/com.teazyou.patch-watcher

	     StartInterval is in seconds; 1800 = 30 min. App updates are rare
	     (days/weeks) so this is deliberately unhurried — change the integer below
	     to re-tune. -->
	<key>Label</key>
	<string>com.teazyou.patch-watcher</string>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/bash</string>
		<string>/Users/teazyou/workspace/scripts/patch-watcher.sh</string>
	</array>
	<key>StartInterval</key>
	<integer>1800</integer>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardOutPath</key>
	<string>/Users/teazyou/workspace/logs/patch_watcher.launchd.out.log</string>
	<key>StandardErrorPath</key>
	<string>/Users/teazyou/workspace/logs/patch_watcher.launchd.err.log</string>
</dict>
</plist>
