# app-configs/aerospace/com.aerospace.display-profile.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.aerospace.display-profile</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/teazyou/workspace/app-configs/aerospace/apply-display-profile.sh</string>
    </array>

    <key>StartInterval</key>
    <integer>5</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/aerospace-display-profile.stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/aerospace-display-profile.stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
